# visa-bot.service
[Unit]
Description=Visa Appointment Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/visa-bot
EnvironmentFile=/opt/visa-bot/.env
ExecStart=/opt/visa-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
